# target_root_script/etc/inittab
::sysinit:/etc/init.d/rcS
console::askfirst:-/bin/sh